# .config/runit/sv/pulseaudio/finish
#!/bin/sh
. $(dirname $0)/common
$pulse --check && exec $pulse -k
